# app/supervisord.conf
; ==============================================================================
; supervisord-Konfiguration
;
; Startet und überwacht einen Prozess im API-Container:
;   1. uvicorn   → FastAPI (REST + MCP + React-UI als Static Files)
;
; Streamlit wurde durch das React-Frontend abgelöst (kein Port 8501 mehr).
;
; Wenn der Prozess abstürzt, wird er automatisch neu gestartet. Logs gehen nach
; stdout/stderr → sichtbar in `docker compose logs api`.
; ==============================================================================

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
; Kein user=root mehr — der Container läuft als non-root (uid 1000, siehe
; Dockerfile USER appuser). supervisord und der fastapi-Kindprozess laufen
; damit als appuser. pidfile nach /tmp, weil dort schreibbar.
pidfile=/tmp/supervisord.pid

[program:fastapi]
command=uvicorn main:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips='*'
directory=/app
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=/app
